# docker/crontab
#minute hour mday month wday who command
* * * * * root cd /var/www/dmi-tcat/capture/stream/ && php controller.php

